# .docker/s6/php-fpm/run
#!/usr/bin/env sh

echo "Starting php-fpm"
exec php-fpm